# scripts/lighthouse_stats.timer
[Unit]
Description=Run lighthouse_stats daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
